# assets/autostart/config.desktop
[Desktop Entry]
Type=Application
Name=Touchpad config
Exec=xbindkeys -p
OnlyShowIn=XFCE
RunHook=0
Hidden=false
